# 2020 - Semestre 1/Análise e Desenvolvimento de Sistemas/1°B/(Equipe JeRso) ProjetoPMIC/qrgen.service
[Unit]
Description=Gerador de codigo qr

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh /home/jerso/qrgen.sh

[Install]
WantedBy=multi-user.target
